# telegram-bot_Version2.service
[Unit]
Description=Telegram Bot Service
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/TelegramBot_test_1
Environment="PATH=/home/YOUR_USERNAME/TelegramBot_test_1/venv/bin"
ExecStart=/home/YOUR_USERNAME/TelegramBot_test_1/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target